# softether.sv.conf
[program:softether]
command=/usr/local/src/softether.sh
redirect_stderr=true
startsecs=0
exitcodes=0
