# softether.sv.conf
[program:softether]
command=vpnserver start
redirect_stderr=true
